# systemd/kx-rdb.service
[Unit]
Description=KX Tickerplant RDB — kxtpstack
After=kx-tp.service
Requires=kx-tp.service

[Service]
Type=simple
EnvironmentFile=-%h/.config/kxtpstack/common.env
# MemoryMax=4G
# CPUQuota=200%
LimitNOFILE=131072
ExecStart=/usr/bin/env bash -lc 'exec "${QBIN:-/opt/kx/q}" "${APP_DIR}/tick/rdb.q" -p "${RDB_PORT:-5011}" -tp "${TP_PORT:-5010}" -hdb "${HDB_PORT:-5012}" ${QOPTS}'
Restart=on-failure
RestartSec=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
