# systemd/kx-tpstack.target
[Unit]
Description=KX Tickerplant Stack 
Wants=kx-tp.service kx-rdb.service kx-hdb.service kx-gw.service
After=kx-tp.service